# deploy/cr1140-app.service
[Unit]
Description=CR1140 baler demo (operator panel)
After=multi-user.target

[Service]
Type=simple
ExecStart=/home/cds-apps/cr1140-baler-demo /dev/input/event1
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
